# launchd/com.aidn.pb-assigner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd job: runs `pb-assigner run` daily at 07:00 local time.

  Install:
    cp launchd/com.aidn.pb-assigner.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.aidn.pb-assigner.plist

  Edit the absolute paths below to match your machine before loading.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aidn.pb-assigner</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/jens-aidn/Documents/Koding/work/PB_assignerV2/.venv/bin/python</string>
        <string>-m</string>
        <string>backend.cli</string>
        <string>run</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jens-aidn/Documents/Koding/work/PB_assignerV2</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/jens-aidn/Library/Logs/pb-assigner.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jens-aidn/Library/Logs/pb-assigner.log</string>

    <!-- Fire daily at 07:00. If the laptop is asleep, launchd will run it on
         next wake; it won't queue multiple missed runs. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
